# ansible/playbooks/group_tasks/proxmox/files/nginx.path
[Unit]
Description=Start nginx when pveproxy-ssl.pem exists

[Path]
PathExists=/etc/pve/nodes/shinx/pveproxy-ssl.pem
Unit=nginx.service

[Install]
WantedBy=multi-user.target