# ansible/playbooks/docker/pancham/files/redis-server.service.d/override.conf
[Unit]
Wants=docker.service
After=network.target docker.service

[Service]
RestartSec=30